# ananicy.service
[Unit]
Description=Ananicy - ANother Auto NICe daemon
After=local-fs.target

[Service]
ExecStart=/usr/bin/ananicy start
Restart=always

[Install]
WantedBy=local-fs.target
